# templates/dnsmasq.service.epp
<% |
    String $after,
    String $memlimit,
| -%>
[Unit]
Description=CF: dnsmasq
After=network.target <%= $after %>
Wants=network.target <%= $after %>

[Service]
ExecStart=/usr/sbin/dnsmasq -k --conf-file=/etc/dnsmasq.conf
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

CPUAccounting=true
CPUShares=100
MemoryAccounting=true
MemoryLimit=<%= $memlimit %>

[Install]
WantedBy=multi-user.target
